# amazoncorretto-17/profile.d/bash.sh
PS1='\u@C17: \W $ '

if [ -d ~/bin ]; then
  export PATH=~/bin:$PATH
fi
