# amazoncorretto-17/profile.d/gpg.sh
#
# Set up the TTY for GPG's user agent and pinentry application.
#
export GPG_TTY=$(tty)
