# amazoncorretto-17/profile.d/java.sh
export JAVA_HOME=/usr/lib/jvm/java-17-amazon-corretto
java -version
